# systemd/memex.service
[Unit]
Description=Memex Knowledge Base Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /usr/lib/node_modules/@wastedcode/memex/dist/index.js serve
Restart=on-failure
RestartSec=5

# Directories (systemd creates these automatically)
RuntimeDirectory=memex
StateDirectory=memex

# Security
AmbientCapabilities=CAP_SYS_ADMIN
NoNewPrivileges=no
CapabilityBoundingSet=CAP_SYS_ADMIN

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
